# luci-app-mosdns/root/etc/hotplug.d/iface/65-mosdns
#!/bin/sh

if [ "$ACTION" = "ifup" ]; then
    /etc/init.d/mosdns restart
fi
