# luci-app-rtorrent/root/etc/init.d/rtorrent
#!/bin/sh /etc/rc.common

START=99
STOP=99

start() {
  screen -dmS rtorrent nice -19 rtorrent -n -o import=/etc/rtorrent.conf
}

boot() {
  start
}

stop() {
  killall rtorrent
}

