# luci-app-ddns-go/root/etc/uci-defaults/ddns-go-permissions
#!/bin/sh

[ -x "/etc/init.d/ddns-go" ] || chmod +x /etc/init.d/ddns-go

exit 0
