# luci-app-proxyxwan/root/etc/uci-defaults/40_luci-app-proxyxwan
touch /etc/config/proxyxwan

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@proxyxwan[-1]
	add ucitrack proxyxwan
	set ucitrack.@proxyxwan[-1].init=proxyxwan
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
