# luci-app-pppwn/root/etc/uci-defaults/luci-pppwn
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@pppwn[-1]
	add ucitrack pppwn
	set ucitrack.@pppwn[-1].init=pppwn
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache

exit 0
