# luci-app-modemband/root/etc/uci-defaults/setupmb.sh
#!/bin/sh
# Copyright 2020-2022 Rafał Wabik (IceG) - From eko.one.pl forum
# MIT License

chmod +x /usr/bin/modemband.sh
chmod +x /usr/bin/loaded.sh
rm -rf /tmp/luci-indexcache
rm -rf /tmp/luci-modulecache/
exit 0
