# luci-app-pptpxwan/root/etc/uci-defaults/40_luci-app-pptpxwan
touch /etc/config/pptpxwan

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@pptpxwan[-1]
	add ucitrack pptpxwan
	set ucitrack.@pptpxwan[-1].init=pptpxwan
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
