# luci-app-music-remote-center/root/etc/uci-defaults/luci-owntone
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@owntone[-1]
	add ucitrack owntone
	set ucitrack.@owntone[-1].init=owntone
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
